# deploy/com.nexus.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexus.gateway</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/npx</string>
        <string>tsx</string>
        <string>packages/cli/src/index.ts</string>
        <string>gateway</string>
        <string>run</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/usr/local/lib/nexus</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>NEXUS_DATA_DIR</key>
        <string>/usr/local/var/lib/nexus</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/nexus/gateway.log</string>

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/nexus/gateway.err.log</string>
</dict>
</plist>
